# deploy/inference-worker.service
[Unit]
Description=Inference Worker (YOLO via Unix socket)
Before=robot-platform.service
After=network.target

[Service]
Type=simple
User=DEPLOY_USER
WorkingDirectory=/opt/robot-platform/inference
ExecStartPre=/bin/rm -f /tmp/inference.sock
ExecStart=/opt/robot-platform/inference/.venv/bin/inference-worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
